# epxcore/src/script/epcore
#!/sbin/runscript


depend() {
	use net.lo
}

# Set START_MODE as "EPS02" or "VGP02", according to start preference

START_MODE="EPS02"
#START_MODE="VGP02"

# Please, don't modify the code after this line ---------------

case $START_MODE in
EPS02)
	VG_NAME=""
	BASEDIR="/opt/eps02"
	;;
VGP02)
	VG_NAME="epcv"
	BASEDIR="/opt/vgp02"
	;;
esac
EP_NAME="epcore"

EP_LOGFILE="$BASEDIR/spool/$EP_NAME.log"
EP_PIDFILE="/var/run/$EP_NAME.pid"
EP_PROPERTIES="$BASEDIR/conf/$EP_NAME.properties"
EP_CORE="$BASEDIR/bin/$EP_NAME"

VG_LOGFILE="$BASEDIR/spool/$VG_NAME.log"
VG_PIDFILE="/var/run/$VG_NAME.pid"
VG_PROPERTIES="$BASEDIR/conf/$VG_NAME.properties"
VG_CORE="$BASEDIR/bin/$VG_NAME"

NULL="/dev/null"
ANSWER="Configured .......: Yes"

start() {
	rm -f $EP_PIDFILE
	rm -f $VG_PIDFILE

	mknod -m 777 /dev/ceusb0 c 198 0
	chmod 777 /dev/ceusb0

	if [ -n "$VG_NAME" ]; then
		mknod -m 777 /dev/ceusb1 c 198 1
		chmod 777 /dev/ceusb1
	fi

	T=0
	until ( ls /dev | grep "ceusb0" > $NULL ); do
		sleep 1s
		T=$(( $T + 1 ))
		[ $T -ge 10 ] && exit 0
	done

	ebegin "Starting ${SVCNAME}"

	setsid $EP_CORE -p $EP_PROPERTIES &> "$EP_LOGFILE" &
	PID=$!
	RET=$?
	echo "$PID" > "$EP_PIDFILE"

	# VG
	if [ -n "$VG_NAME" ]; then
		setsid $VG_CORE -p $VG_PROPERTIES &> "$VG_LOGFILE" &
		PID=$!
		RET=$?
		echo "$PID" > "$VG_PIDFILE"
	fi
        
	eend $RET
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	if [ -e $EP_PIDFILE ] ; then
		kill $(< "$EP_PIDFILE")
		rm "$EP_PIDFILE"
	fi
	if [ -e $VG_PIDFILE ] ; then
		kill $(< "$VG_PIDFILE")
		rm "$VG_PIDFILE"
	fi
	eend 0
}

